# crontab
0 5 * * * make -C /elpa > /var/log/elpa-cron.log
